# systemd/llmos-web.service
[Unit]
Description=LLM-OS Web UI Server
After=ollama.service network-online.target
# Wants (not Requires): the web UI must still start and bind port 8080 even if
# Ollama is slow to come up or failed — the UI reports Ollama status itself.
Wants=ollama.service

[Service]
Type=simple
User=llmos
Group=llmos
ExecStart=/usr/local/bin/llmos --web --port 8080
Restart=on-failure
RestartSec=5
Environment="HOME=/home/llmos"
Environment="LLMOS_MODEL=llama3.2"
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
